# backend/thinkphp/supervisor/ai_dev_worker.conf
[program:ai-dev-code-worker]
directory=/path/to/thinkphp-project
command=php think queue:work --queue ai_dev_code --tries 2 --sleep 2 --timeout 3600
autostart=true
autorestart=true
startsecs=3
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/ai-dev-code-worker.log
stopwaitsecs=30
